# test/crontabs/mycronfile
* * * * * docker ps > /var/log/dockerps.log 2>&1
* * * * * /myscripts/script.sh
